# etc/systemd/user/stretch-reminder.timer
[Unit]
Description=Timer for Stretch Reminder

[Timer]
OnCalendar=*-*-* *:45:00
Persistent=true

[Install]
WantedBy=timers.target
