# etc/crontab/jobs
15,45   *   *   *   *   notify.sh "Please Stand Up" >/dev/null 2>&1
