# etc/crontab/jobs
15,45   *   *   *   *   /home/jacob/.local/bin/notify.sh "Please Stand Up" >/dev/null 2>&1
